# debian/cron.d/quantum-dhcp-agent-netns-cleanup
# vim: set filetype=crontab:
# Periodically cleans Quantum's network namespaces on behalf of the Quantum
# DHCP agent.
30 * * * * quantum if [ -x /usr/bin/quantum-netns-cleanup ] ; then /usr/bin/quantum-netns-cleanup --config-file=/etc/quantum/quantum.conf --config-file=/etc/quantum/dhcp_agent.ini >/dev/null 2>&1; fi
